# otto.service
[Unit]
Description=OttO IoT Server
Documentation=https://github.com/rustyeddy/otto
After=network.target

[Service]
Type=simple
User=otto
Group=otto
WorkingDirectory=/opt/otto
ExecStart=/opt/otto/bin/otto serve
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/otto/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otto

[Install]
WantedBy=multi-user.target
